# camera-app-foreground.service
[Unit]
Description=Raspberry Pi Camera App (Foreground)
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/bin/bash -c 'cd /home/pi/camera_app_raspi2 && /usr/bin/python3 camera_app.py'
StandardInput=tty-force
StandardOutput=inherit
StandardError=inherit
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
Environment=TERM=linux
User=pi
Group=pi
Restart=on-failure
RestartSec=10

# Wait for system to be ready
ExecStartPre=/bin/sleep 15

[Install]
WantedBy=multi-user.target 