# camera-app.service
[Unit]
Description=Raspberry Pi Camera App
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=miura
Group=miura
WorkingDirectory=/home/pi/camera_app_raspi2
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin
ExecStart=/usr/bin/python3 /home/pi/camera_app_raspi2/camera_app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Wait for camera to be ready
ExecStartPre=/bin/sleep 10

[Install]
WantedBy=multi-user.target 